# templates/grafana.service
[Unit]
Description=Grafana
Documentation=https://grafana.com/docs/grafana/latest/
After=network.target

[Service]
Type=simple
User=ubuntu
ExecStart=/home/ubuntu/grafana-11.1.3/bin/grafana-server --config=/home/ubuntu/grafana-11.1.3/conf/defaults.ini
Restart=always
WorkingDirectory=/home/ubuntu/grafana-11.1.3
StandardOutput=file:/var/log/grafana.log
StandardError=file:/var/log/grafana-error.log

[Install]
WantedBy=multi-user.target
